# infra/systemd/mission-control-prod.service
[Unit]
Description=Mission Control PROD (Docker Compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/kuba/repos/mission-control
EnvironmentFile=/etc/mission-control/prod.env
ExecStartPre=/usr/bin/bash -lc '/home/kuba/repos/mission-control/infra/scripts/run-api-migrations.sh /home/kuba/repos/mission-control/infra/prod/docker-compose.prod.yml /etc/mission-control/prod.env'
ExecStart=/usr/bin/docker compose -f /home/kuba/repos/mission-control/infra/prod/docker-compose.prod.yml --env-file /etc/mission-control/prod.env up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f /home/kuba/repos/mission-control/infra/prod/docker-compose.prod.yml --env-file /etc/mission-control/prod.env down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
