# infra/mission-control.service
[Unit]
Description=Mission Control — Supervisor System Dashboard
After=network.target

[Service]
Type=simple
User=kuba
WorkingDirectory=/home/kuba/.openclaw/workspace-james/repos/mission-control/apps/web
EnvironmentFile=/home/kuba/mission-control/mission-control.env
ExecStart=/usr/bin/env npm start -- -p 3100
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
