# infra/mission-control.service
[Unit]
Description=Mission Control — Supervisor System Dashboard
After=network.target

[Service]
Type=simple
User=kuba
WorkingDirectory=/home/kuba/.openclaw/workspace-james/repos/mission-control/apps/web
Environment=NODE_ENV=production
Environment=SUPERVISOR_SYSTEM_PATH=/home/kuba/.openclaw/SUPERVISOR_SYSTEM
ExecStart=/usr/bin/env npm start -- -p 3100
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
